# scripts/com.crypto.discovery.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crypto.discovery</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/leandrosilva/Projects/personal/cryptocurrencies-board/scripts/run-discovery.sh</string>
    </array>

    <!-- Run monthly on day 1 at 18:00 UTC (adjust Hour if your local timezone differs from UTC) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Do not run on load -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/leandrosilva/Projects/personal/cryptocurrencies-board</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/leandrosilva/Projects/personal/cryptocurrencies-board/logs/discovery.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/leandrosilva/Projects/personal/cryptocurrencies-board/logs/discovery.err.log</string>

    <!-- Keep alive settings -->
    <key>KeepAlive</key>
    <false/>

    <!-- Nice level (lower priority) -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- Longer timeout for discovery (30 minutes) -->
    <key>TimeOut</key>
    <integer>1800</integer>
</dict>
</plist>
